# ops/systemd/sogecon-web.service
[Unit]
Description=Sogecon Web (Next.js standalone)
After=network.target

[Service]
Type=simple
# 운영 계정/그룹을 환경에 맞게 조정하세요.
User=sogecon
Group=sogecon
# `web-deploy.sh`가 관리하는 릴리스 심볼릭 링크 경로
# 기본 경로는 /opt/sogecon/web/current 입니다.
# 레포 내부에 릴리스를 두려면 ops/web-deploy.sh 의 RELEASE_BASE 를 변경하고,
# 여기 WorkingDirectory 도 동일한 경로로 맞추세요 (예: /srv/sogecon-app/.releases/web/current).
WorkingDirectory=/opt/sogecon/web/current
Environment=NODE_ENV=production
Environment=PORT=3000
# 주의: NEXT_PUBLIC_* 변수는 빌드 타임에 고정됩니다. 런타임에서 변경해도 반영되지 않습니다.
# 빌드 시점에 .env.local/.env.production 등을 통해 고정하세요.
# Next.js standalone 진입점(릴리스 루트 기준)
ExecStart=/usr/bin/node apps/web/server.js
Restart=always
RestartSec=3
KillSignal=SIGINT
TimeoutStopSec=15
SyslogIdentifier=sogecon-web
NoNewPrivileges=true
ProtectSystem=full
ProtectHome=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
